# bootstrap/ansible/roles/networking/files/firewall/nftables.service.d/override.conf
[Service]
ProtectSystem=strict
# Rely on the fact that nftables.service is a oneshot service
# that does not clear rules on "stop" (it has no stop command).
ExecStopPost=-/etc/nftables/persisted-sets.sh
