# bootstrap/desktop-environments/profile.d/70-dev.sh
# export EDITOR='emacsclient -c -a emacs'
# export EDITOR='gvim -f'
export EDITOR='vim'
# export PATH=$HOME/.config/emacs/bin:$PATH
