# bootstrap/desktop-environments/profile.d/71-git_askpass.sh
if [[ "$XDG_CURRENT_DESKTOP" =~ .*KDE ]]; then
    export GIT_ASKPASS=/usr/bin/ksshaskpass
fi
