# bootstrap/desktop-environments/profile.d/71-pagers.sh
export PAGER="less"
export RI="-f ansi"
export LESS="-Ri -x4"
